# kuu-bot.service
[Unit]
Description=Kuu Telegram Bot - Anonymous Chat Service
Documentation=https://github.com/royalashu4u/kuu
After=network.target
Wants=network.target

[Service]
Type=notify
User=kuu-bot
Group=kuu-bot
WorkingDirectory=/root/bot
Environment=PATH=/root/bot/.venv/bin
EnvironmentFile=/root/bot/.env
ExecStart=/root/bot/.venv/bin/python app.py
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kuu-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/bot/data
ReadWritePaths=/root/bot/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Process management
PIDFile=/root/bot/kuu-bot.pid
NotifyAccess=main

[Install]
WantedBy=multi-user.target
